# etc/systemd/system/mautic_segments_update.timer
[Unit]
Description=test

[Timer]
OnUnitInactiveSec=30s
OnBootSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
